# service/rsyslog/run
#!/bin/bash
set -e

rm -f /var/run/rsyslogd.pid

exec rsyslogd -n
